# client.service.example
[Unit]
Description=Pipecat Assistant Frontend Client
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/pipecat-assistant/client
Environment=PORT=2774
Environment=NODE_ENV=production
ExecStart=/path/to/npm run start  # which npm
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
